# bot/minecraft-bot.service
[Unit]
Description=Minecraft Discord Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/MC_server_logbot/bot
ExecStart=/home/ubuntu/MC_server_logbot/.venv/bin/python main.py
Restart=always
RestartSec=5
EnvironmentFile=/home/ubuntu/MC_server_logbot/bot/.env

[Install]
WantedBy=multi-user.target
